# smf/manifests/napi.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Networking API SMF manifest

	Copyright (c) 2012, Joyent, Inc. All rights reserved.
 -->

<service_bundle type="manifest" name="napi">
<service name="smartdc/application/napi" type="service" version="1">

<dependency name="network" grouping="require_all" restart_on="error" type="service">
	<service_fmri value="svc:/network/physical" />
</dependency>
<dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
	<service_fmri value="svc:/system/filesystem/local" />
</dependency>

<exec_method type="method" name="start" exec="/opt/smartdc/napi/bin/napi &amp;" timeout_seconds="30" />
<exec_method type="method" name="stop" exec=":kill" timeout_seconds="30" />

<instance name="default" enabled="true" />

<stability value='Unstable' />

<template>
	<common_name><loctext xml:lang="C">Joyent Networking API</loctext></common_name>
</template>

</service>
</service_bundle>
